# assets/etc/supervisord/php.conf
[program:php]
command = /usr/sbin/php-fpm7
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0
autorestart = false
startentries = 0
